# deployment/vacation-dashboard.service
[Unit]
Description=Vacation Dashboard Web Application
After=network.target

[Service]
Type=exec
User=vacation-dashboard
Group=vacation-dashboard
WorkingDirectory=/opt/vacation-dashboard
Environment=PATH=/opt/vacation-dashboard/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/vacation-dashboard
EnvironmentFile=-/etc/vacation-dashboard/.env
ExecStartPre=/bin/mkdir -p /var/run/vacation-dashboard
ExecStartPre=/bin/chown vacation-dashboard:vacation-dashboard /var/run/vacation-dashboard
ExecStart=/opt/vacation-dashboard/venv/bin/python -u app.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Безпека
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vacation-dashboard /var/log/vacation-dashboard /var/run/vacation-dashboard /tmp

[Install]
WantedBy=multi-user.target